# packaging/lxfu-face.service
[Unit]
Description=LXFU Face Verification Service
After=network.target

[Service]
Type=simple
ExecStart=@PREFIX@/bin/lxfu_face_service
Restart=on-failure
# Uncomment to run under a dedicated account
# User=lxfu
# Group=video

[Install]
WantedBy=multi-user.target
